# crystalbudget.service
[Unit]
Description=CrystalBudget Flask Application
After=network.target

[Service]
Type=exec
User=root
Group=root
WorkingDirectory=/home/admin/crystall-budget
Environment="SECRET_KEY=0xm5hy67Txo5taHWK_CuWy72wSDFkgN0f9Rc2L7c_ZhdYFMJTYu8IoXjoW6YQy0x"
Environment="BUDGET_DB=/home/admin/crystall-budget/budget.db"
Environment="HTTPS_MODE=false"
Environment="LOG_LEVEL=INFO"
ExecStart=/home/admin/crystall-budget/venv/bin/gunicorn \
    --bind 127.0.0.1:5000 \
    --workers 2 \
    --timeout 120 \
    --keep-alive 5 \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    --preload \
    --access-logfile /home/admin/crystall-budget/access.log \
    --error-logfile /home/admin/crystall-budget/error.log \
    --log-level info \
    app:app

# Автоматический перезапуск при падении
Restart=always
RestartSec=10

# Безопасность
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/admin/crystall-budget /home/admin/crystall-budget
PrivateTmp=true

# Лимиты ресурсов
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target