# crystall-web.service
[Unit]
Description=Crystall Budget Next.js Web
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=crystall
Group=crystall
WorkingDirectory=/data/crystall-budget/web
ExecStart=/usr/bin/node node_modules/next/dist/bin/next start -p 3000
Environment=NODE_ENV=production
Restart=always
RestartSec=5
RuntimeMaxMemory=300M
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target