# crystalbudget.service
[Unit]
Description=CrystalBudget Personal Finance Application
After=network.target
Wants=network.target

[Service]
Type=notify
User=admin
Group=admin
WorkingDirectory=/opt/crystall-budget

# Environment
Environment="FLASK_ENV=production"
Environment="LOG_LEVEL=INFO"
Environment="SECRET_KEY=0xm5hy67Txo5taHWK_CuWy72wSDFkgN0f9Rc2L7c_ZhdYFMJTYu8IoXjoW6QQy0x"
Environment="BUDGET_DB=sqlite:////var/lib/crystalbudget/budget.db"
Environment="PYTHONPATH=/opt/crystall-budget"
Environment="HTTPS_MODE=true"
Environment="TELEGRAM_BOT_TOKEN=test-token"
Environment="MODAL_SYSTEM_ENABLED=true"

# Ensure directories exist
ExecStartPre=/bin/mkdir -p /var/lib/crystalbudget
ExecStartPre=/bin/chown admin:admin /var/lib/crystalbudget

# Production server
ExecStart=/opt/crystall-budget/.venv/bin/gunicorn \
  --chdir /opt/crystall-budget \
  --workers=3 --threads=2 --timeout=60 \
  --bind 127.0.0.1:5000 \
  app:app

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crystalbudget

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
ReadWritePaths=/var/lib/crystalbudget

# Resource limits
LimitNOFILE=1024
LimitNPROC=512

[Install]
WantedBy=multi-user.target